# deploy/launchd/com.gombwe.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- gombwe — macOS LaunchAgent (installer template).
     Native autostart for macOS hosts. One hop: launchd -> gombwe (no pm2).

     This is a LaunchAgent (per-user GUI session), NOT a LaunchDaemon, ON PURPOSE:
     the policy scanner spawns `claude -p`, whose Claude SUBSCRIPTION token lives
     in the user's LOGIN KEYCHAIN. Only a process inside the user's login (Aqua)
     session can read that keychain — a system LaunchDaemon runs in a different
     security session and CANNOT reach it, even after the user logs in. So a
     LaunchDaemon breaks AI classification unless you switch to a billed
     ANTHROPIC_API_KEY. LaunchAgent keeps it on the subscription, free.
       KeepAlive = crash recovery; RunAtLoad = start at login.
       MUST run `start --headless` — plain `start` opens a REPL that dies with no
       TTY and triggers a KeepAlive restart loop.

     Trade-off: an Agent starts at login, not at boot. On a FileVault host that's
     moot — a human must unlock at the console each boot anyway, and that unlock
     logs them in, which fires the Agent. If a user truly needs pre-login start
     AND AI classification, that requires an API key + a LaunchDaemon instead.

     Install (no sudo needed for a user agent):
       cp com.gombwe.daemon.plist ~/Library/LaunchAgents/
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.gombwe.daemon.plist -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.gombwe.daemon</string>
    <key>ProgramArguments</key>
    <array>
      <string>@NODE_BIN@</string>
      <string>@GOMBWE_DIST@/index.js</string>
      <string>start</string>
      <string>--headless</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>WorkingDirectory</key>
    <string>@GOMBWE_REPO@</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>NODE_ENV</key><string>production</string>
      <key>HOME</key><string>@HOME@</string>
      <key>PATH</key><string>@HOME@/.local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StandardOutPath</key><string>@HOME@/.claude-gombwe/gombwe.out.log</string>
    <key>StandardErrorPath</key><string>@HOME@/.claude-gombwe/gombwe.err.log</string>
  </dict>
</plist>
